# mount_autorun/Library/LaunchAgents/st.tomfro.unmount_autorun.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
 <dict>
  <key>KeepAlive</key>
  <dict>
   <key>SuccessfulExit</key>
   <false/>
  </dict>
  <key>Label</key>
  <string>st.tomfro.unmount_autorun</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>PATH=$PATH:~/bin ~/bin/unmount_autorun_monitor</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>/tmp/ua_err.log</string>
  <key>StandardOutPath</key>
  <string>/tmp/ua_out.log</string>
 </dict>
</plist>
